# net-irc/atheme/files/atheme.initd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-irc/atheme/files/atheme.initd,v 1.1 2006/04/11 21:09:53 swegener Exp $

depend() {
	need net
}

start() {
	ebegin "Starting atheme"
	start-stop-daemon --start --quiet --chuid atheme --exec /usr/bin/atheme &>/dev/null
	eend $? "Failed to start atheme"
}

stop() {
	ebegin "Stopping atheme"
	start-stop-daemon --stop --quiet --pidfile /var/lib/atheme/var/atheme.pid
	eend $? "Failed to stop atheme"
}
